# toInstall/Leap_Motion_SDK_Linux_2.3.1/LeapDeveloperKit_2.3.1+31549_linux/temp/etc/init/leapd.conf
# leapd - Leap Motion daemon

description	"Leap Motion daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec leapd
